# debian/snapbackend.cron.weekly
#!/bin/sh -e
#

# Run a check against the list data to make sure it remains valid over
# long period of time... (until we use a proper BATCH command to handle
# these INSERTs)
#
# Make sure that this script is active only on the backend computer(s)
# that run the snappagelist
#
if test -s /usr/bin/snapconfig && /usr/bin/snapconfig snapserver backends | grep -q snappagelist
then
    snapbackend --action links::cleanuplinks
fi

# vim: ts=4 sw=4 et
